# etc/supervisord.d/rsyslog_conf.ini
[program:rsyslog-conf]
command=consul-template -consul localhost:8500 -wait 5s -template "/etc/consul-templates/rsyslog_targets.conf.ctmpl:/etc/rsyslog.d/targets.conf:supervisorctl restart rsyslog"
redirect_stderr=true
autostart=true
startsecs=15
stdout_logfile=/var/log/supervisor/%(program_name)s.log
